# supervisord.conf
[supervisord]
nodaemon=true

[program:whatsmeow_bot]
command=/app/bot
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/bot.err.log
stdout_logfile=/var/log/supervisor/bot.out.log
